# roles/app/templates/gunicorn.socket
[Unit]
Description=gunicorn-{{ instance_name }} socket

[Socket]
ListenStream={{ gunicorn_socket }}

[Install]
WantedBy=sockets.target
